# scripts/local-macos/templates/launchd/local.distillery-rescore.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>local.distillery-rescore</string>

    <key>ProgramArguments</key>
    <array>
        <string>__HOME__/.distillery/rescore.sh</string>
    </array>

    <!-- Daily at 04:15 local time. -->
    <key>StartCalendarInterval</key>
    <dict>
        <key>Hour</key>
        <integer>4</integer>
        <key>Minute</key>
        <integer>15</integer>
    </dict>

    <key>StartCalendarIntervalRunOnMissed</key>
    <true/>

    <key>StandardOutPath</key>
    <string>__HOME__/.distillery/rescore.out.log</string>
    <key>StandardErrorPath</key>
    <string>__HOME__/.distillery/rescore.err.log</string>

    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/opt/homebrew/bin:/usr/bin:/bin:/usr/sbin:/sbin</string>
    </dict>
</dict>
</plist>
